# scripts/app.service
[Unit]
Description=Zen Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zen
ExecStart=/opt/zen/app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
